# .config/autostart/plasma-bspwm.sh.desktop
[Desktop Entry]
Exec=/home/kuzy000/bin/plasma-bspwm.sh
Icon=dialog-scripts
Name=plasma-bspwm.sh
Path=
Type=Application
X-KDE-AutostartScript=true
